# deploy/ai-pyramid/pet-album.service
[Unit]
Description=Pet Album Server (ai-pyramid)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/smart-pet-camera/src/ai-pyramid
EnvironmentFile=-/opt/smart-pet-camera/.env

Environment=PET_ALBUM_PHOTOS_DIR=/opt/smart-pet-camera/src/ai-pyramid/data/photos
Environment=PET_ALBUM_DB_PATH=/opt/smart-pet-camera/src/ai-pyramid/data/pet-album.db
Environment=PET_ALBUM_VLM_URL=http://localhost:8000
Environment=PET_ALBUM_VLM_MODEL=AXERA-TECH/Qwen3-VL-2B-Instruct-GPTQ-Int4-C256-P3584-CTX4095

ExecStartPre=/bin/mkdir -p ${PET_ALBUM_PHOTOS_DIR}
ExecStart=/opt/smart-pet-camera/src/ai-pyramid/target/release/pet-album \
  --photos-dir "${PET_ALBUM_PHOTOS_DIR}" \
  --db-path "${PET_ALBUM_DB_PATH}" \
  --vlm-url "${PET_ALBUM_VLM_URL}" \
  --vlm-model "${PET_ALBUM_VLM_MODEL}"

Restart=on-failure
RestartSec=3

[Install]
WantedBy=ai-pyramid.target
