# deploy/rdk-x5/pet-camera.target
[Unit]
Description=Smart Pet Camera Stack (rdk-x5)
After=network-online.target
Wants=network-online.target
Wants=pet-camera-capture.service
Wants=pet-camera-detector.service
Wants=pet-camera-monitor.service
Wants=pet-camera-streaming.service
Wants=comic-sync.service

[Install]
WantedBy=multi-user.target
